# sanjaya/sanjaya.service
[Unit]
Description=SANJAYA Shipment Intelligence API
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/sanjaya
ExecStart=/home/ec2-user/sanjaya/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
